# sign-tool.service
[Unit]
Description=Sign Tool Web Service
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/sign-tool
ExecStart=/opt/sign-tool/.venv/bin/python -m sign_tool web --host 127.0.0.1 --port 2087
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
